# dotfiles/.config/profile.d/000-nix
#!/usr/bin/env bash

# Ensure we have configured Nix.
sourceExists '/nix/var/nix/profiles/default/etc/profile.d/nix-daemon.sh'
