# dotfiles/.config/profile.d/030-nix
#!/usr/bin/env bash

sourceExists '/nix/var/nix/profiles/default/etc/profile.d/nix-daemon.sh'
