# dotfiles/.config/profile.d/i3
#!/usr/bin/env bash

export TERMINAL='konsole'
