# dotfiles/.config/profile.d/intellij-ce
#!/bin/bash

prependPathsExists ~/Applications/IntelliJ\ IDEA\ 14\ CE.app/Contents/MacOS
