# files/home/.config/profile.d/020-snap
#!/usr/bin/env bash

if [[ -d /snap/bin ]]; then
  export PATH="$PATH:/snap/bin"
fi
