# files/home/.config/profile.d/gnu-sed
#!/usr/bin/env bash

gnuSedBase="/usr/local/opt/gnu-sed/libexec"

if [[ -d $gnuSedBase ]]; then
  PATH="$gnuSedBase/gnubin:$PATH"
  MANPATH="$gnuSedBase/gnuman:$MANPATH"
fi
