# files/home/.config/profile.d/golang
#!/bin/bash

export GOPATH=$HOME/.go
PATH=$GOPATH/bin:$PATH
