# files/home/.config/profile.d/golang
#!/usr/bin/env bash

export GOPATH=$HOME/.go
PATH=$GOPATH/bin:$PATH
