# files/home/.config/profile.d/hunspell
#!/usr/bin/env bash

export DICTIONARY=en_GB
export DICPATH=~/.nix-profile/share/hunspell
