# files/home/.config/profile.d/lean
#!/usr/bin/env bash

# Configure the Lean version manager (ala rustup).
sourceExists ~/.elan/env
