# files/home/.config/profile.d/rbenv
#!/bin/bash

if true; then
  if [[ -d ~/.rbenv ]]; then
    PATH=~/.rbenv/bin:${PATH}
    eval "$(rbenv init -)"
  fi
fi

# vim: set filetype=sh
