# files/home/.config/profile.d/ripgrep
#!/bin/bash

export RIPGREP_CONFIG_PATH=~/.config/ripgrep/config
