# files/home/.config/profile.d/rust
#!/usr/bin/env bash

sourceExists ~/.cargo/env
