# files/home/.config/profile.d/rvm
#!/usr/bin/env bash

if false; then
  if [[ -d ~/.rvm ]]; then
    export PATH="$PATH:$HOME/.rvm/bin"
    sourceExists ~/.rvm/scripts/rvm
  fi
fi
