# files/home/.config/profile.d/yarn
#!/usr/bin/env bash

PATH=~/.yarn/bin:$PATH
